# templates/pulp/pulp.service.erb
[Unit]
Description=<%= @description %>
After=docker.service <%= @after_units.join(' ') %>
Require=docker.service <%= @required_units.join(' ') %>

[Service]
ExecStart=/usr/bin/podman start --conmon-pidfile=%t/<%= @container %>.pid --cidfile=%t/<%= @container %>.ctr-id --cgroups=no-conmon <%= @container %>
ExecStop=/usr/bin/podman stop -t 10 <%= @container %>
Restart=always
RestartSec=10s
# Trust podman to stop the container
KillMode=None
Type=forking
TimeoutStartSec=120
TimeoutStopSec=15
PIDFile=%t/<%= @container %>.pid

[Install]
WantedBy=multi-user.target
