# templates/jenkins-agent.service.erb
[Unit]
Description=ROS build farm jenkins agent
After=network.target

[Service]
EnvironmentFile=/etc/default/<%= @service_name %>
Type=simple
ExecStart=/bin/sh -c "/usr/bin/java $JAVA_ARGS -jar $SWARM_CLIENT_JAR \
-master $JENKINS_URL -username $USERNAME -password $PASSWORD \
-name $NAME -description \"$DESCRIPTION\" -mode $MODE -executors $EXECUTORS \
-fsroot $FSROOT -labels \"$LABELS\" $JENKINS_ARGS"

User=<%= @username %>
Restart=always

[Install]
WantedBy=multi-user.target
